# services/acolock_back.socket
[Unit]
Description=ACoLock Backend socket

[Socket]
ListenStream=/run/acolock/socket

[Install]
WantedBy=sockets.target
